# rootfs/overlay/etc/systemd/system/NetworkManager.service.d/10-boot-resilience.conf
# Bound NetworkManager's startup so a flapping / absent USB-Ethernet dongle can't
# wedge the boot.
#
# Debian's NetworkManager.service ships TimeoutStartSec=10min, and NM is
# Type=dbus, so systemd doesn't consider the boot's default target "reached"
# until NM claims org.freedesktop.NetworkManager on the bus. On felix the only
# NIC is a USB-Ethernet dongle (r8152) that can enumerate slowly (~130s) or flap
# or be absent (see project_gs201_usb_host_mode_port); when NM's startup wedges
# on it, that 10-minute ceiling keeps the boot "in progress" the entire time.
# getty.target itself is already up by then (the *-issue generators don't block),
# but while the boot is unfinished systemd keeps spraying per-service status onto
# the ttySAC0 console every second, which drowns the serial login prompt — so the
# box looks hung with no way in over UART. Observed: an 11+ minute "hang".
#
# Fix: cap the start to 90s and auto-restart. Worst case the box reaches a clean,
# reachable multi-user state in ~90s instead of 10 minutes, and NM keeps retrying
# so it comes up on its own once the dongle finally enumerates. 90s is far more
# than NM needs to claim its bus name (a healthy boot here finishes in ~13s);
# device bring-up happens *after* startup, so this never trips on a merely-slow
# dongle.
[Unit]
# A flapping dongle can make NM fail its start repeatedly; don't let the default
# start-rate-limit trip and leave NM permanently dead — keep retrying.
StartLimitIntervalSec=0

[Service]
TimeoutStartSec=90
Restart=on-failure
RestartSec=10
