# rootfs/overlay/etc/systemd/system/NetworkManager.service.d/10-bounded-stop.conf
[Service]
# See ../../system.conf.d/10-bounded-shutdown.conf — these three were the
# observed offenders that hung shutdown for 10+ minutes on felix.
TimeoutStopSec=10s
