# rootfs/overlay/etc/systemd/system/kmsconvt@.service.d/10-wait-network.conf
[Unit]
# Delay the login prompt until DHCP has settled so /etc/issue's \4 escape
# renders a real IPv4 address. Wants= (not Requires=) means the console still
# comes up if the wait-online service times out (default 30s).
Wants=NetworkManager-wait-online.service
After=NetworkManager-wait-online.service
