# rootfs/overlay/etc/systemd/system/kmsconvt@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/kmscon "--vt=%I" --seats=seat0 --no-switchvt --login -- /sbin/agetty - xterm-256color
