# rootfs/overlay/etc/systemd/system/battery-issue.service
[Unit]
Description=Update battery readout in the login banner (/etc/issue)
ConditionPathExists=/etc/issue

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/battery-issue
