# rootfs/overlay/etc/systemd/system/custom.service
[Unit]
Description=Custom helpers

[Service]
Type=oneshot
ExecStart=/usr/bin/customservice

[Install]
WantedBy=multi-user.target
