# rootfs/overlay/etc/systemd/system/mark-slot-successful.service
[Unit]
Description=Mark current Pixel A/B slot successful (devinfo)
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Mark the running slot (from androidboot.slot_suffix) successful in the devinfo
# partition: successful, max retries, clear unbootable. Stops the bootloader's
# retry counter from exhausting after a few boots and forcing fastboot mode.
# devinfo-ONLY (does NOT touch the UFS boot LUN), so the retry reset keeps
# working even if the boot-LUN sysfs node can't be located.
# Source: github.com/junkyard-computing/pixel-bootctl (vendored at
# tools/pixel-bootctl and cross-compiled to aarch64 by the build).
ExecStart=/usr/local/bin/pixel-bootctl mark-successful
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
