# rootfs/overlay/etc/systemd/system/usb-gadget.service
[Unit]
Description=Configure CDC-NCM + CDC-ACM USB gadget (felix mainline)
# Run AFTER getty / login prompts are rendered. The script writes verbose
# [USB-GADGET] / [DWC3-DBG] traffic to /dev/kmsg, which intermixes with
# agetty's terminal-size probe sequence on ttySAC0; running early made
# login-prompt rendering intermittent. Putting this after getty.target
# means agetty's probe completes (and the prompt is drawn) before the
# gadget service starts and writes to kmsg.
After=getty.target sys-kernel-config.mount
Requires=sys-kernel-config.mount
# No-op unless the Type-C port landed device-role and dwc3 exposed a UDC.
# With the dongle attached the port is data=host, there is no UDC, and this
# service does not run -- so it cannot disturb the ethernet+charging path.
# Both modes are required, and this condition is what keeps them disjoint.
# Also started on UDC hot-plug via 99-usb-gadget.rules, because the cable is
# moved between dongle and PC at runtime, not just at boot.
ConditionPathExistsGlob=/sys/class/udc/*

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/usb_gadget
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
